# home/setup/.config/autostart/oem-setup.desktop
[Desktop Entry]
Type=Application
Name=OEM Setup
Exec=/usr/local/bin/oem-setup.sh
X-GNOME-Autostart-enabled=true
NoDisplay=true
